# Website/myscript.service
[Unit]
Description=My Script Service
After=multi-user.target

[Service]
Type=idle
<<<<<<< HEAD
ExecStart=sudo /usr/bin/python3 /var/www/html/VMMaster.py
=======
ExecStart=/usr/bin/sudo /usr/bin/python3 /var/www/html/VMMaster.py
>>>>>>> origin/master

[Install]
WantedBy=multi-user.target
